# launchd/com.micahbos.observatory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.micahbos.observatory</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/micahbos/Developer/observatory/launchd/start.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>OBSERVATORY_PORT</key>
    <string>3099</string>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/micahbos/Developer/observatory/server</string>
  <key>StandardOutPath</key>
  <string>/Users/micahbos/.observatory/observatory.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/micahbos/.observatory/observatory.err</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>
</dict>
</plist>
